# udev-rules/99-hdmicheck.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/etc/lightdm/display_hdmi.sh"
